# rsyslog.socket
[Unit]
Description=Syslog Socket

[Socket]
ListenDatagram=/dev/log

[Install]
WantedBy=sockets.target
