# rsyslog.service.in
[Unit]
Description=System Logging Service

[Service]
ExecStart=@sbindir@/rsyslogd -n
Sockets=syslog.socket
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=syslog.service
