# rsyslog.service.in
[Unit]
Description=System Logging Service

[Service]
ExecStart=@sbindir@/rsyslogd -n -c5

[Install]
WantedBy=multi-user.target
Also=rsyslog.socket
