# rsyslog.service.in
[Unit]
Description=System Logging Service

[Service]
ExecStart=@sbindir@/rsyslogd -n -c 4

[Install]
WantedBy=multi-user.target
